# etc/init/project_name.conf
start on runlevel [2345]
stop on runlevel [!2345]

expect daemon
respawn

exec /usr/local/bin/uwsgi --ini {{ uwsgi_inifile }}
